# systemd/rk_video_ai.service
# systemd service unit: RK3568 双路视频三模型推理服务
# 安装路径: /etc/systemd/system/rk_video_ai.service

[Unit]
Description=RK3568 dual-video three-model inference service
# 等待网络就绪、MQTT Broker 和 3A 服务完成初始化后再启动本服务
After=network-online.target mosquitto.service rkaiq_3A.service
Wants=network-online.target mosquitto.service rkaiq_3A.service
# /userdata 必须为独立挂载分区，防止误用根文件系统（根分区可能为只读或空间不足）
RequiresMountsFor=/userdata
ConditionPathIsMountPoint=/userdata

[Service]
Type=simple
# 工作目录设为可执行文件所在目录，相对路径（config/model 等）基于此解析
WorkingDirectory=/userdata/rk3568_inspection_terminal/bin
# 传感器配置文件路径，通过环境变量注入，程序启动时读取
Environment=EDGE_SENSOR_CONFIG=/userdata/rk3568_inspection_terminal/config/sensors.ini
# 可选的环境变量覆盖文件，前缀 "-" 表示文件不存在时不报错
EnvironmentFile=-/etc/default/rk3568-inspection
# 启动前等待 10 秒，确保网络、MQTT Broker 等依赖服务完全就绪
ExecStartPre=/bin/sleep 10
ExecStart=/userdata/rk3568_inspection_terminal/bin/rk3568_inspection_terminal
# 仅在异常退出时重启（非零退出码或信号终止），正常退出不重启
Restart=on-failure
# 以下值不要使用行尾注释，旧版 systemd 会把注释当作值的一部分。
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM
LimitNOFILE=4096
LimitRTPRIO=99
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal

[Install]
# 允许 systemctl enable；是否开机启动仍由部署者决定。
WantedBy=multi-user.target
